# floodgate.service
[Unit]
Description=FloodGate XDP Traffic Scrubber
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/floodgate -i vlan50 -s 10
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
